# deploy/systemd/hc-meta-sync.service
[Unit]
Description=HC Agency Meta background sync
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/hc-agency-dashboard/current
EnvironmentFile=-/etc/hc-agency-dashboard.env
ExecStart=/usr/bin/flock -n /tmp/hc-meta-sync.lock /usr/bin/env node scripts/sync.js
Nice=5
IOSchedulingClass=best-effort
IOSchedulingPriority=6
